# script/root/sv/log/run
#!/bin/sh
mkdir -p /var/log/opendatanow
svlogd -ttt /var/log/opendatanow
